# deploy/supervisor/horizon.conf.example
[program:dotsheet-horizon]
process_name=%(program_name)s
command=/usr/bin/php /var/www/dotsheet/artisan horizon
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/dotsheet/storage/logs/horizon.log
stopwaitsecs=3600
numprocs=1
environment=APP_ENV="production",HOME="/var/www/dotsheet"
